# etc/system/systemd/mashup-timer.timer
[Unit]
Description=Meshup Countdown Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=60sec

[Install]
WantedBy=timers.target
