# etc/system/systemd/mashup.timer
[Unit]
Description=Run mashup script every second

[Timer]
OnUnitActiveSec=1s
AccuracySec=1s
Unit=mashup.service

[Install]
WantedBy=timers.target
